# debian/ipu7-camera-dynamic.service
[Unit]
Description=IPU7 Camera Dynamic Bridge (on-demand LED)
After=graphical-session.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/ipu7-camera
ExecStart=/usr/sbin/ipu7-camera-dynamic
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
